# config/systemd/active/domtech.service
# =============================================================================
# DOMTECH - SYSTEMD SERVICE UNIT
# =============================================================================
[Unit]
Description=Dominguez Tech Solutions - Technology Consulting and Development Services
Documentation=https://domingueztechsolutions.com
After=network.target mysql.service nginx.service
Wants=network-online.target
Requires=mysql.service

[Service]
Type=simple
User=admin
Group=admin

# Site-specific working directory
WorkingDirectory=/home/admin/deno-genesis/sites/domingueztechsolutions-com

# Standard DenoGenesis environment variables
Environment=DENO_ENV=production
Environment=PORT=3000
Environment=DB_HOST=localhost
Environment=DB_USER=webadmin
Environment=DB_PASSWORD=Password123!
Environment=DB_NAME=universal_db
Environment=SITE_KEY=domtech

# Deno execution with comprehensive permissions for DenoGenesis framework
ExecStart=/home/admin/.deno/bin/deno run \
  --allow-net \
  --allow-read \
  --allow-write \
  --allow-env \
  --allow-run \
  --unstable \
  --lock=deno.lock \
  --cached-only \
  main.ts

# Failure recovery configuration
Restart=on-failure
RestartSec=10
StartLimitBurst=3
StartLimitInterval=60

# Resource limits for production stability
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening measures
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

# File system access permissions
ReadWritePaths=/home/admin/deno-genesis/sites/domingueztechsolutions-com
ReadWritePaths=/tmp
ReadWritePaths=/var/log

# Process management configuration
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
TimeoutStartSec=60

# Logging configuration for systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=denogenesis-domtech

[Install]
WantedBy=multi-user.target
Alias=domtech.service